# docker/cron/crontab
0 3 * * * /usr/local/bin/php /var/www/html/artisan waybill:check >> /var/log/cron.log 2>&1
